# orchids.service.sample
[Unit]
Description=Orchids climate control
After=syslog.target

[Service]
Type=simple
WorkingDirectory=/root/savemyorchids
ExecStart=/root/savemyorchids/.venv/bin/python /root/savemyorchids/run.py
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
